# supervisor.conf
; Supervisor program for the BOOM Kafka-consumer micro-service.
;
; baselayer (skyportal/baselayer: tools/setup_services.py -> copy_supervisor_configs)
; reads this file VERBATIM when it is present in the service directory, instead of
; generating a single-instance program. That lets us run multiple consumers.
;
; All instances share the same Kafka consumer group.id (set in main.py), so Kafka
; auto-distributes the topic partitions across them via consumer.subscribe(). This
; gives both throughput headroom for LSST-volume streams and resilience / gap-free
; deploys (one instance can restart while the others keep ingesting). Re-delivery
; on rebalance is safe: duplicate candidates are rejected by the IntegrityError
; guard in main.py, and offsets are committed only after a successful write
; (consumer.commit(..., asynchronous=True)).
;
; numprocs is the number of consumers to run. Useful parallelism is bounded by the
; partition count on the ZTF_alerts_results / LSST_alerts_results topics (extra
; consumers beyond that sit idle) and by Postgres write headroom — raise it in
; lockstep with upstream topic partitioning and DB capacity, not blindly. Default
; 2 gives resilience without assuming any extra partitions.
;
; (Making numprocs config-driven would require a baselayer change to template this
; file during copy_supervisor_configs, which currently reads it verbatim.)
[program:boom]
command=/usr/bin/env python services/boom/main.py %(ENV_FLAGS)s
environment=PYTHONPATH=".",PYTHONUNBUFFERED="1"
numprocs=2
process_name=%(program_name)s_%(process_num)s
stdout_logfile=log/boom_service_%(process_num)s.log
redirect_stderr=true
